# scratch/single-tenant-bundle-target-host/systemd/strategy-validator.service
# Backend-only single-tenant Docker service template.
# Copy to /etc/systemd/system/strategy-validator.service and adjust paths/image.
[Unit]
Description=Strategy Validator single-tenant backend API
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
Environment=STRATEGY_VALIDATOR_IMAGE=strategy-validator-api:local
Environment=STRATEGY_VALIDATOR_HOST_PORT=8000
EnvironmentFile=/etc/strategy-validator/deployment.env
ExecStartPre=/usr/bin/docker run --rm --user 0:0 \
  --read-only \
  --tmpfs /tmp:rw,nosuid,nodev,noexec,size=64m \
  --cap-drop ALL \
  --security-opt no-new-privileges:true \
  -v /etc/strategy-validator:/deployment-env:ro \
  ${STRATEGY_VALIDATOR_IMAGE} \
  strategy-validator-deployment-env-check /deployment-env/deployment.env --require-valid --json
ExecStartPre=-/usr/bin/docker rm -f strategy-validator-api
ExecStart=/usr/bin/docker run --rm --name strategy-validator-api \
  --env-file /etc/strategy-validator/deployment.env \
  -p 127.0.0.1:${STRATEGY_VALIDATOR_HOST_PORT}:8000 \
  --read-only \
  --tmpfs /tmp:rw,nosuid,nodev,noexec,size=64m \
  --cap-drop ALL \
  --security-opt no-new-privileges:true \
  -v strategy-validator-ledger:/var/lib/strategy-validator \
  -v strategy-validator-backups:/var/backups/strategy-validator \
  ${STRATEGY_VALIDATOR_IMAGE}
ExecStop=-/usr/bin/docker stop strategy-validator-api
Restart=on-failure
RestartSec=5s
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
